# roles/pg_dump/templates/pg_dump.timer.j2
[Unit]
Description=Dump PostgreSQL database and upload to S3 - daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
