# roles/pg_dump/templates/pg_dump.service.j2
[Unit]
Description=Dump PostgreSQL database and upload to S3

[Service]
Type=oneshot
ExecStart=/usr/bin/ansible-playbook {{ pg_dump_home }}/do_pg_dump.yaml
Environment=HOME=%h
